# circadianlight.service
[Unit]
Description=Circadian Light

[Service]
ExecStart=%h/.cargo/bin/circadianlight daemon

[Install]
WantedBy=graphical.target
